# saslauthd.ini
[program:saslauthd]
command="service saslauthd start"